# systemd/user/openclaw-voice-transcriber-cleanup.timer
[Unit]
Description=Run OpenClaw Voice Transcriber inbound audio cleanup

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
